# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:vncserver]
command=/usr/bin/vncserver :1 -geometry 1280x720 -depth 24 -nolisten tcp
autorestart=true
autostart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/vncserver.log
priority=998

[program:xfce4]
command=/usr/bin/startxfce4
autorestart=false
autostart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/xfce4.log
environment=DISPLAY=:1
priority=999

[program:lyrion-sync]
command=python /app/src/main.py
autorestart=true
autostart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/lyrion-sync.log
environment=DISPLAY=:1,HOME=/root
priority=999
